# gmail-idle.service
[Unit]
Description=Gmail IDLE to MQTT Bridge
After=network.target

[Service]
Type=simple
User=watanabeyousuke
WorkingDirectory=/home/watanabeyousuke/MassageApp
ExecStart=/home/watanabeyousuke/venv-gmail/bin/python /home/watanabeyousuke/MassageApp/gmail_idle_to_mqtt.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target